# gcon.service
[Unit]
Description=gcon server

[Service]
ExecStart=/bin/zsh -c '/usr/bin/gcon-server'

[Install]
WantedBy=default.target
